# bin/init-script
#! /bin/sh
### BEGIN INIT INFO
# Provides:          kvm-machines
# Required-Start:    $kvm $network
# Required-Stop:     
# Default-Start:     2 3 4 5
# Default-Stop:      S 0 1 6
# Short-Description: KVM virtual machines init script
# Description:       This script starts the kvm virtual machines
### END INIT INFO

. /lib/lsb/init-functions

KVM_DEFS=/etc/default/kvm-manager

if [ ! -f $KVM_DEFS ] 
then
	echo "ERROR: $KVM_DEFS does not exist."
	exit 3
fi

. $KVM_DEFS

if [ "$KVM_MANAGER_ROOT" = "" ]
then
	echo "ERROR: The KVM_MANAGER_ROOT variable is not defined $KVM_DEFS"
	exit 3
fi

KVM_MANAGER=$KVM_MANAGER_ROOT/bin/kvm-manager
if [ ! -f $KVM_MANAGER ]
then
	echo "ERROR: The KVM_MANAGER_ROOT variable is not defined properly in $KVM_DEFS"
	echo "It should be the directory where you installed kvm-manager"
	exit 3
fi

case "$1" in
  start)
	# Not starting all machines at the same time
	# the most urgent services are up and running first
	for m in $ENABLED_MACHINES
	do
		$KVM_MANAGER $m start
		sleep 5
	done
	exit 0
	;;
  stop)
	# Stopping the least critical machines first
	for m in `$KVM_MANAGER_ROOT/bin/reverse-list "$ENABLED_MACHINES"`
	do
		$KVM_MANAGER $m stop
	done
	exit 0
	;;
  restart)
	# Restart the machines one by one
	# Better in terms of uptime than stopping all machines first
	# and starting them all again then
	for m in $ENABLED_MACHINES
	do
		$KVM_MANAGER $m stop
		$KVM_MANAGER $m start
	done
	exit 0
	;;
  reload-iptables)
	/etc/network/if-up.d/iptables
	for m in $ENABLED_MACHINES
	do
		$KVM_MANAGER $m disable-maintenance
		$KVM_MANAGER $m load-iptables
	done
	exit 0
	;;
  enable-maintenance)
	for m in $ENABLED_MACHINES
	do
		$KVM_MANAGER $m enable-maintenance
	done
	$KVM_MANAGER_ROOT/bin/kvm-info
	exit 0
	;;
  disable-maintenance)
	for m in $ENABLED_MACHINES
	do
		$KVM_MANAGER $m disable-maintenance
	done
	$KVM_MANAGER_ROOT/bin/kvm-info
	exit 0
	;;
  *)
	echo "Usage: $0 {start|stop|reload-iptables|enable-maintenance|disable-maintenance|restart}" >&2
	exit 3
	;;
esac
